# files/iop-cvemap-download.path
[Unit]
Description=Path unit for cvemap.xml manual file watcher

[Path]
PathChanged=/var/lib/foreman/cvemap.xml
PathModified=/var/lib/foreman/cvemap.xml

[Install]
WantedBy=multi-user.target
